# packages/ducatuscore-node/ducatuscore-node.service
[Unit]
Description=Ducatuscore Node
After=network.target

[Service]
User=ducatuscore
Group=ducatuscore
ExecStart=/usr/bin/node /home/ducatuscore/ducatuscore/packages/ducatuscore-node/build/src/server.js

[Install]
WantedBy=multi-user.target
